# CSharp_090122_Cron/crontab
30  12       * * 1-2,4-5 Cinani
15  12       * * 3       Menza
30  11,14    * * 7,6     Svaca doma
30  6        * * *       Snidane doma
0   19       * * *       Vecere doma